# cambiarPropietarioPines.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cambiarPropietarioPines.sh

[Install]
WantedBy=multi-user.target
