# powerOff-E135.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/powerOff-E135 &

[Install]
WantedBy=multi-user.target
