# send_email.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/send_email.sh &

[Install]
WantedBy=multi-user.target
